# root/etc/services.d/tools/run
#!/usr/bin/with-contenv bash

cd /app
exec s6-setuidgid simon python3 run.py
